# django-asgi.service
[Unit]
Description=Django ASGI Server for WebSocket connections
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/EmployeeAttendance
Environment=DJANGO_SETTINGS_MODULE=attendance_system.settings
Environment=ENVIRONMENT=production
Environment=PYTHONPATH=/var/www/EmployeeAttendance
ExecStart=/var/www/EmployeeAttendance/venv/bin/daphne -b 127.0.0.1 -p 8001 attendance_system.asgi:application
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=django-asgi

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/EmployeeAttendance
ReadWritePaths=/var/log
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
